# etc/systemd/system/getty@tty1.service.d/override.conf
[Service]
ExecStart=
ExecStart=-/usr/bin/journalctl -u minecraft -f
StandardInput=tty
StandardOutput=tty
